# dist/xyz.openbmc_project.Logging.service
[Unit]
Description=Phosphor Log Manager

[Service]
ExecStart=/usr/bin/phosphor-log-manager
Restart=always
Type=dbus
BusName=xyz.openbmc_project.Logging

[Install]
WantedBy=multi-user.target
Alias=dbus-xyz.openbmc_project.Logging.service
