# extensions/amd-event-logging/service/amd-log-manager.service
[Unit]
Description=AMD Event Logging Manager (AEL)
After=xyz.openbmc_project.Logging.service
Wants=xyz.openbmc_project.Logging.service

[Service]
ExecStart=/usr/bin/amd-log-manager
Restart=always
RestartSec=5
Type=simple

[Install]
WantedBy=multi-user.target
